# systemd_units/update_scrapper.timer
[Unit]
Description=Run RPi Update Scrapper Script daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
